# systemctl_services/hitos_cameras.service
# sudo ln -s /home/arvc/ros2_ws/src/hitos_setup/systemctl_services/hitos_cameras.service /etc/systemd/system/hitos_cameras.service
[Unit]
Description=HITOS Cameras - Multiespectral FLIR+Basler acquisition
After=network.target hitos_ptp.service hitos_ptp_sync.service hitos_sensors.service
Wants=network.target hitos_ptp.service hitos_ptp_sync.service hitos_sensors.service

[Service]
Type=simple
User=arvc
# The session folder is generated here and shared with hitos_sync.service via
# /tmp/hitos_session.env so cameras and buffer compositor use the SAME dataset
# session. hitos_sync has PartOf=hitos_cameras.service: restarting the cameras
# (new session) automatically restarts the sync side too. The session timestamp
# is generated in ExecStartPre (BEFORE the unit goes active) so hitos_sync, which
# is ordered After= this unit, always reads the fresh value. Fixes the race where
# sync read /tmp/hitos_session.env before ExecStart had written it (sync landed a
# session behind, e.g. data in mult_test_session / a stale calib_<date>).
ExecStartPre=/bin/bash -c "echo HITOS_SESSION=$(date +'%%d-%%m-%%Y_%%Hh%%Mm') > /tmp/hitos_session.env"
ExecStart=/bin/bash -c "\
    source /home/arvc/ros2_ws/install/setup.bash && \
    source /home/arvc/ros2_ws/install/hitos_setup/share/hitos_setup/scripts/hitos_setup.sh && \
    _hitos_export_ip && \
    export SPINNAKER_GENTL64_CTI=/opt/spinnaker/lib/spinnaker-gentl/Spinnaker_GenTL.cti && \
    HITOS_CALIB=0; HITOS_CALIB_VISIBLE_HZ=4; [ -f /tmp/hitos_mode.env ] && . /tmp/hitos_mode.env; \
    RATE=1; [ $HITOS_CALIB = 1 ] && RATE=$HITOS_CALIB_VISIBLE_HZ; \
    . /tmp/hitos_session.env; SESSION=$HITOS_SESSION; \
    ros2 launch multiespectral_acquire cameras_only.launch.py \
        session_folder:=$SESSION output_frame_rate:=$RATE"
KillSignal=SIGINT
TimeoutStopSec=45
Restart=on-failure
RestartSec=65

[Install]
WantedBy=multi-user.target
