# systemctl_services/hitos_log_cleanup.timer
[Unit]
Description=HITOS — daily log cleanup timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
