# systemctl_services/logrotate.timer.d/hourly.conf
# HITOS drop-in: run logrotate HOURLY instead of daily.
#
# Why: /var/log/syslog is bounded only by logrotate's `maxsize 100M` rule
# (/etc/logrotate.d/rsyslog). `maxsize` is evaluated ONLY when logrotate runs.
# With the stock daily timer, a high-rate logger (e.g. the u-blox GPS "End of
# file" spin-loop, ~340 MB/h) can balloon syslog to several GB between runs —
# observed: 7.4 GB in one cycle. Running hourly makes `maxsize 100M` effective,
# capping the live syslog to a few hundred MB and total syslog footprint to
# ~0.5 GB even during a sustained log flood.
#
# Time-based rules (daily/weekly in other logrotate configs) are unaffected —
# they still rotate on their own schedule; only `maxsize` rules trigger early.
#
# Install:
#   sudo mkdir -p /etc/systemd/system/logrotate.timer.d
#   sudo ln -sf /home/arvc/ros2_ws/src/hitos_setup/systemctl_services/logrotate.timer.d/hourly.conf \
#               /etc/systemd/system/logrotate.timer.d/hourly.conf
#   sudo systemctl daemon-reload
#   sudo systemctl restart logrotate.timer
# Verify:  systemctl show logrotate.timer -p TimersCalendar   # -> OnCalendar=*-*-* *:00:00
#          systemctl list-timers logrotate.timer

[Timer]
# Empty assignment first resets the inherited OnCalendar=daily, then set hourly.
OnCalendar=
OnCalendar=hourly
